# metatrader/supervisord.conf
[supervisord]
nodaemon=true

[program:startup]
command=/root/startup.sh
autorestart=true
priority=1000